# rootfs/etc/supervisor.d/bind.ini

[program:bind9]
command=/usr/sbin/named -c /etc/bind/named.conf -u named -f -g
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true
