# rootfs/etc/supervisor.d/samba4.ini

[program:samba]
command=/usr/sbin/samba --configfile=/srv/etc/samba/smb.conf --interactive
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true
